# packages/vtuner-apps/unit/vtuner.service
[Unit]
Description=Virtual Tuner client/server

[Service]
EnvironmentFile=/etc/vtunerc.conf
ExecStart=/usr/sbin/vtuner.sh
Restart=always

[Install]
WantedBy=multi-user.target
